# systemd/email-sync.timer
[Unit]
Description=Timer to sync emails

[Timer]
OnUnitInactiveSec=1m

[Install]
WantedBy=timers.target
