# systemd/user.config/backup.timer
[Unit]
Description=Timer to backup files

[Timer]
OnCalendar=*-*-* 0:23:00
Persistent=true

[Install]
WantedBy=timers.target
