# systemd/user.config/backblaze.timer
[Unit]
Description=Timer to upload backups to backblaze

[Timer]
OnCalendar=*-*-* 19:23:00
Persistent=true

[Install]
WantedBy=timers.target
