# systemd/user.config/newsboat.timer
[Unit]
Description=Timer to download RSS feeds

[Timer]
OnBootSec=5m
OnUnitInactiveSec=30m

[Install]
WantedBy=timers.target
